# dotenv/zsh/profile.d/00_path.zsh
PATH=
_add_path /bin
_add_path /usr/bin
_add_path /usr/local/bin
_add_path $HOME/bin
_add_path /sbin
_add_path /usr/sbin
_add_path /usr/local/sbin
