# dotenv/zsh/profile.d/goenv.zsh
_add_path $HOME/.goenv/bin
_add_path $HOME/.goenv/packages/bin
