# dotenv/zsh/profile.d/tig.zsh
_add_path $HOME/.dotfiles/var/plugins/git/tig/bin
